# docker-portal-supervisor.conf
[program:docker-portal]
directory=/vagrant
command=node app.js
numprocs=1
autostart=true
autorestart=true
